# dotfiles/launchagents/local.zcode.usagebar.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>local.zcode.usagebar</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/goldenwave/Applications/ZCodeUsageBar/zcode-usage-bar/ZCodeUsageBar.app/Contents/MacOS/ZCodeUsageBar</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>LimitLoadToSessionType</key>
  <string>Aqua</string>
  <key>ProcessType</key>
  <string>Interactive</string>
  <key>StandardOutPath</key>
  <string>/tmp/zcode-usagebar.launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/zcode-usagebar.launchd.err.log</string>
</dict>
</plist>
